# Proxmox/sh/zfs-check.timer
[Unit]
Description= check ZFS status
[Timer]
OnBootSec=5min
OnCalendar=*-*-* 09:03:00
[Install]
WantedBy=timers.target